# portage-tmpfs
#!/sbin/openrc-run
# nic@boet.cc
# Populate a tmpfs-mounted Portage repo tree at /var/db/repos from a snapshot tarball.
# Provides commands: setup remove snapshot (plus normal start/stop).

description="Populate a tmpfs-mounted Portage repo tree from a snapshot tarball"
description_snapshot="Archive the current repo tree to a compressed tarball"
description_setup="Create snapshot, purge repo tree, add tmpfs fstab entry"
description_remove="Unmount tmpfs, restore on-disk repo, remove fstab entry"

extra_started_commands="snapshot"
extra_stopped_commands="setup remove"

readonly PORTAGE_ROOT="/var/db/repos"
readonly TARBALL="/var/db/repos-snapshot.tar.gz"
readonly SENTINEL="${PORTAGE_ROOT}/tmpfs-not-mounted"

# Marker used to find/remove the exact fstab line this service owns
readonly FSTAB_MARK="## ${SVCNAME} ##"

depend() {
    before nfs
}

# -------- helpers --------

# Print an error message and return failure.
die() {
    eerror "   $*"
    return 1
}

# Return 0 if the given command exists in PATH.
have_cmd() {
    command -v "$1" >/dev/null 2>&1
}

# Return 0 if the given path is an active mountpoint.
# Prefer mountpoint(1); otherwise fall back to /proc/mounts.
is_mounted() {
    if have_cmd mountpoint; then
        mountpoint -q -- "$1"
    fi

    # POSIX-ish fallback
    # shellcheck disable=SC2002
    awk -v mp="$1" '$2 == mp { found=1 } END { exit(found ? 0 : 1) }' /proc/mounts 2>/dev/null
}

# Return 0 if the gentoo repo metadata timestamp file exists.
has_gentoo_repo_marker() {
    [ -f "${PORTAGE_ROOT}/gentoo/metadata/timestamp" ]
}

# Remove and recreate a directory safely.
# Refuse to operate on empty, root, non-absolute, or symlinked paths.
safe_purge_dir() {
    local _d="$1"

    case "${_d}" in
        ""|"/")
            die "refusing to purge unsafe path '${_d}'"
            return 1
            ;;
        /*)
            ;;
        *)
            die "refusing to purge non-absolute path '${_d}'"
            return 1
            ;;
    esac

    if [ -L "${_d}" ]; then
        die "refusing to purge symlink '${_d}'"
        return 1
    fi

    rm -rf -- "${_d}" || return 1
    mkdir -p -- "${_d}" || return 1
    return 0
}

# Return 0 if /etc/fstab contains a line with our service marker.
fstab_has_our_entry() {
    grep -F -- "${FSTAB_MARK}" /etc/fstab >/dev/null 2>&1
}

# Append a tmpfs fstab entry for PORTAGE_ROOT tagged with our marker.
# No-op if the entry already exists.
fstab_add_entry() {
    if fstab_has_our_entry; then
        ewarn "   fstab already contains entry for ${SVCNAME}"
        return 0
    fi

    echo "none  ${PORTAGE_ROOT}  tmpfs  defaults,size=2G  0  0  ${FSTAB_MARK}" >> /etc/fstab || return 1
    return 0
}

# Delete any fstab line containing our service marker.
fstab_remove_entry() {
    if ! fstab_has_our_entry; then
        return 0
    fi

    # Delete only lines containing our marker
    sed -i "\|${FSTAB_MARK}|d" /etc/fstab || return 1
    return 0
}

# Return 0 if the tarball exists and was modified within the last 24 hours.
tarball_is_recent() {
    [ -f "${TARBALL}" ] || return 1

    # Use find -print -quit and check for any output
    if find "${TARBALL}" -mmin -1440 -print -quit 2>/dev/null | grep -q .; then
        return 0
    fi

    return 1
}

# -------- commands --------

# Archive PORTAGE_ROOT into a gzipped tarball with absolute paths.
# Skips creation if a recent (<24h) tarball already exists.
snapshot() {
    if tarball_is_recent; then
        einfo "   skipping snapshot creation, recent archive exists"
        return 0
    fi

    if [ -f "${SENTINEL}" ]; then
        die "refusing to create snapshot, sentinel exists: ${SENTINEL}"
        return 1
    fi

    if [ ! -d "${PORTAGE_ROOT}" ]; then
        die "portage root does not exist: ${PORTAGE_ROOT}"
        return 1
    fi

    einfo "   creating portage snapshot at ${TARBALL}"

    # Archive the whole PORTAGE_ROOT tree using absolute paths.
    # -c create, -z gzip, -p preserve perms, -f file
    # Use --absolute-names so the archive stores leading /
    if ! tar -czpf "${TARBALL}" --absolute-names "${PORTAGE_ROOT}" >/dev/null 2>&1; then
        die "snapshot failed (tar returned non-zero)"
        return 1
    fi

    return 0
}

# Validate and extract the snapshot tarball back to its original absolute paths.
restore() {
    if [ ! -f "${TARBALL}" ]; then
        die "snapshot tarball missing: ${TARBALL}"
        return 1
    fi

    # Validate tarball before extracting: must contain expected paths only
    if ! tar -tzf "${TARBALL}" 2>/dev/null | grep -q "^${PORTAGE_ROOT}"; then
        die "tarball does not contain expected ${PORTAGE_ROOT} tree"
        return 1
    fi

    einfo "   restoring ${TARBALL} to /"

    # Extract using -C / and --absolute-names to honor absolute paths.
    if ! tar -xzpf "${TARBALL}" -C / --absolute-names >/dev/null 2>&1; then
        die "restore failed (tar returned non-zero)"
        return 1
    fi

    return 0
}

# First-run deployment: snapshot the on-disk repo, purge PORTAGE_ROOT,
# place a sentinel file, and add a tmpfs fstab entry.
setup() {
    ebegin "Deploying ${SVCNAME}"

    if fstab_has_our_entry; then
        ewarn "   already setup: fstab entry exists for ${SVCNAME}"
        eend 0
        return 0
    fi

    if [ -f "${SENTINEL}" ]; then
        ewarn "   already setup: sentinel exists: ${SENTINEL}"
        eend 0
        return 0
    fi

    if is_mounted "${PORTAGE_ROOT}"; then
        die "mountpoint already mounted at ${PORTAGE_ROOT} (refusing to modify)"
        eend 1
        return 1
    fi

    if ! has_gentoo_repo_marker; then
        die "existing gentoo repo not found under ${PORTAGE_ROOT} (missing metadata/timestamp)"
        eend 1
        return 1
    fi

    snapshot || { eend 1; return 1; }

    einfo "   purging ${PORTAGE_ROOT}"
    safe_purge_dir "${PORTAGE_ROOT}" || { eend 1; return 1; }

    einfo "   creating sentinel ${SENTINEL}"
    touch -- "${SENTINEL}" || { eend 1; return 1; }

    einfo "   adding tmpfs mount to /etc/fstab"
    fstab_add_entry || { eend 1; return 1; }

    einfo "setup complete, please start service"
    eend 0
    return 0
}

# Tear down: unmount tmpfs, remove sentinel and fstab entry,
# restore the on-disk repo from tarball, then delete the tarball.
remove() {
    ebegin "Removing ${SVCNAME}"

    if is_mounted "${PORTAGE_ROOT}"; then
        einfo "   unmounting ${PORTAGE_ROOT}"
        if ! umount -fl -- "${PORTAGE_ROOT}" >/dev/null 2>&1; then
            die "failed to unmount ${PORTAGE_ROOT}"
            eend 1
            return 1
        fi

        if is_mounted "${PORTAGE_ROOT}"; then
            die "env not stable, ${PORTAGE_ROOT} still mounted"
            eend 1
            return 1
        fi
    fi

    if [ -f "${SENTINEL}" ]; then
        rm -f -- "${SENTINEL}" || { eend 1; return 1; }
    fi

    fstab_remove_entry || { eend 1; return 1; }

    einfo "   purging ${PORTAGE_ROOT}"
    safe_purge_dir "${PORTAGE_ROOT}" || { eend 1; return 1; }

    restore || { eend 1; return 1; }

    if [ -f "${TARBALL}" ]; then
        rm -f -- "${TARBALL}" || { eend 1; return 1; }
    fi

    einfo "complete"
    eend 0
    return 0
}

# Mount the tmpfs at PORTAGE_ROOT and restore the snapshot if empty.
start() {
    ebegin "Starting ${SVCNAME}"

    if is_mounted "${PORTAGE_ROOT}"; then
        ewarn "   ${PORTAGE_ROOT} is mounted; unmounting first"
        if ! umount -fl -- "${PORTAGE_ROOT}" >/dev/null 2>&1; then
            die "failed to unmount ${PORTAGE_ROOT}"
            eend 1
            return 1
        fi
    fi

    if [ ! -f "${SENTINEL}" ]; then
        die "please run setup first (missing sentinel: ${SENTINEL})"
        eend 1
        return 1
    fi

    # mount(8) uses /etc/fstab entry we added
    if ! mount -- "${PORTAGE_ROOT}" >/dev/null 2>&1; then
        die "failed to mount ${PORTAGE_ROOT}"
        eend 1
        return 1
    fi

    # If tmpfs is empty, restore snapshot
    if ! has_gentoo_repo_marker; then
        restore || { eend 1; return 1; }
    fi

    eend 0
    return 0
}

# Snapshot the live repo tree and unmount the tmpfs.
stop() {
    ebegin "Stopping ${SVCNAME}"

    if is_mounted "${PORTAGE_ROOT}"; then
        if has_gentoo_repo_marker; then
            snapshot || ewarn "   failed to archive current tmpfs contents"
        fi

        if ! umount -fl -- "${PORTAGE_ROOT}" >/dev/null 2>&1; then
            ewarn "   failed to unmount tmpfs"
            # do not hard-fail stop; report via eend
            eend 1
            return 1
        fi
    else
        ewarn "   ${PORTAGE_ROOT} not mounted"
    fi

    eend 0
    return 0
}
